# .docker/supervisord-goip.conf
[program:goip_cron]
command=/usr/local/goip/goipcron_x64 /usr/local/goip/inc/config.inc.php
numprocs=1
autostart=true
