# .docker/supervisord-goip.conf
[program:goipcron]
command = /usr/local/goip/goipcron_x64 /usr/local/goip/inc/config.inc.php
startretries = 0
autorestart = false
